# src/rove_launch_handler/scripts/launch_script.service
[Unit]
Descript=Robot launch script

[Service]
ExecStart=/usr/bin/launch_script.sh

[Install]
WantedBy=multi-user.target